# mctomqtt.service
[Unit]
Description=MeshCore to MQTT Relay
After=time-sync.target
Wants=time-sync.target

[Service]
# Run as the user who installed npm packages, NOT root
# Replace <username> with your actual username
User=<username>
WorkingDirectory=</path/to/meshcoretomqtt>
ExecStart=update_and_start.sh
KillMode=process
Restart=on-failure
Type=exec

[Install]
WantedBy=multi-user.target
